# clojure-lein/templates/clojure.conf.erb
#!upstart
description "clojure lein ring server"
author      "OpDemand"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

script
    export HOME="<%= scope.lookupvar("clojure::params::home") %>"
    cd <%= scope.lookupvar("clojure:params:repository_path") %> 
    exec sudo -u <%= scope.lookupvar("clojure::params::username") %> -i lein deps
    exec sudo -u <%= scope.lookupvar("clojure::params::username") %> -i lein ring server 8080
end script
