# clojure/templates/clojure.conf.erb
#!upstart
description "clojure lein ring server"
author      "OpDemand"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

script
    export HOME=/home/ubuntu/repo
    exec sudo -u <%= scope.lookupvar("clojure::params::username") %> -i lein ring server-headless <%= scope.lookupvar("clojure::params::port") %> 2>&1 >> /var/log/nodejs.log
end script
